# services.d/icinga2-cli-daemon/run
#!/usr/bin/execlineb -P
icingacli director daemon
