# services.d/nginx/run
#!/usr/bin/execlineb -P
# nginx -g "daemon off;" 2>/proc/1/fd/2 1>/proc/1/fd/1
nginx -g "daemon off;"
